# ops/systemd/tradeos-ai.service
[Unit]
Description=TradeOS AI Next.js application
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/var/www/tradeos-ai
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=/etc/tradeos-ai/tradeos-ai.env
ExecStart=/usr/bin/npm run start:prod
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
